# src-ui/service/osdp-MONITOR.service
[Unit]
Description=OSDP Monitor Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=no
User=root
WorkingDirectory=/opt/osdp-conformance/run/MON
ExecStart=/opt/osdp-conformance/bin/open-osdp

[Install]
WantedBy=multi-user.target



